# deploy/systemd/nft-miner-api.service.example
[Unit]
Description=NFT Miner API
Documentation=file:/opt/portfolio-demo-app/docs/DEPLOY_CHECKLIST.md
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=nftminer
Group=nftminer
WorkingDirectory=/opt/portfolio-demo-app
EnvironmentFile=/etc/nft-miner-api.env
Environment=NODE_ENV=production
Environment=API_HOST=127.0.0.1
Environment=API_PORT=3101
ExecStart=/usr/bin/node server/src/index.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/portfolio-demo-app
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
